# kyg/v2x/systemd/v2x-recorder.service
[Unit]
Description=V2X Circular Recorder (FFmpeg 10s segments)
After=network-online.target
[Service]
Type=simple
# === RTSP 입력 === (예: rtsp://user:pass@CAM_IP/stream)
ExecStart=/usr/bin/ffmpeg -nostdin -y -rtsp_transport tcp -i "RTSP_URL_HERE" \
  -c copy -f segment -segment_time 10 -strftime 1 -reset_timestamps 1 \
  /var/rec/seg_%s.mp4
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
