# kyg/v2x/systemd/v2x-watcher.service
[Unit]
Description=V2X Event Watcher (trigger broadcast + clip)
After=v2x-recorder.service network-online.target
Wants=v2x-recorder.service

[Service]
WorkingDirectory=/opt/v2x
ExecStart=/usr/bin/python3 /opt/v2x/watch_and_broadcast.py
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
